# dist/sshReverseProxy.service
[Unit]
Description=A layer-7 reverse proxy for the SSH/SFTP protocol
After=nss-lookup.target

[Service]
Restart=always
RestartSec=30
EnvironmentFile=/etc/default/sshreverseproxy
ExecStartPre=-/bin/mkdir -p /var/run/sshreverseproxy/
ExecStart=/usr/bin/sshreverseproxy --config /etc/sshreverseproxy/sshreverseproxy.conf --loglevel $LOGLEVEL daemon
LimitNOFILE=65635
